# macos/ttyd-wrapper.plist
<?xml version="1.0" encoding="UTF-8"?>
<!DOCTYPE plist PUBLIC "-//Apple//DTD PLIST 1.0//EN" "http://www.apple.com/DTDs/PropertyList-1.0.dtd">
<!-- ttyd-wrapper LaunchAgent (template - rendered by install-service.sh) -->
<plist version="1.0">
<dict>
    <key>Label</key>
    <string>com.pawprint0706.ttyd-wrapper</string>

    <key>ProgramArguments</key>
    <array>
        <string>__TTYD__</string>
        <string>--writable</string>
        <string>-p</string>
        <string>__PORT__</string>
        <string>-I</string>
        <string>__INDEX__</string>
        <string>--cwd</string>
        <string>__HOME__</string>
        <string>zsh</string>
        <string>-l</string>
    </array>

    <key>RunAtLoad</key>
    <true/>

    <key>KeepAlive</key>
    <dict>
        <key>SuccessfulExit</key>
        <false/>
    </dict>

    <key>StandardOutPath</key>
    <string>__LOG__</string>
    <key>StandardErrorPath</key>
    <string>__LOG__</string>
</dict>
</plist>
