# linux/ttyd-wrapper.service
# ttyd-wrapper systemd user unit (template - rendered by install-service.sh)
[Unit]
Description=ttyd Web Terminal
After=network.target

[Service]
ExecStart=__TTYD__ __ARGS__
Restart=on-failure
RestartSec=3

[Install]
WantedBy=default.target
